# install/clash.timer
[Unit]
Description=Restart Clash Service daily

[Timer]
OnCalendar=*-*-* 8,20:00:00
Persistent=true

[Install]
WantedBy=timers.target